# update-mmdb-systemd/update-mmdb.timer
[Unit]
Description=Get latest MMDB

[Timer]
OnCalendar=Mon *-*-* 00:00:00

[Install]
WantedBy=timers.target
